# conf.d/lxde.conf
[program:xterm]
command=startlxde
autorestart=true
